# cookbooks/lan-vpc-route/files/vpc-route.service
[Unit]
# Static routes from a LAN host to the AWS VPC and the tailnet CGNAT range,
# next-hopped at the pro-router LXC (192.168.1.60 = home tailnet subnet router).
#
# Why this is needed: the LAN default gateway (RTX1210, 192.168.1.253) carries a
# static route for these CIDRs -> 192.168.1.60 (home-monitor rtx_static_route.aws
# / .tailscale), but the RTX does NOT hairpin-forward to a next-hop that sits on
# the same LAN segment as the source, so traffic sent via the default gateway is
# silently dropped. Hosts must therefore hold the route directly. DHCP clients
# learn it via DHCP option 121 (home-monitor rtx_dhcp_scope); static-IP hosts
# (PVE LXCs) get it from this unit.
#
# CIDRs/next-hop are fleet constants (see home-monitor contracts/devices.json
# pro-router = .60, var.vpc_cidr = 10.33.128.0/18, var.tailscale_cgnat_cidr =
# 100.64.0.0/10). Update here if those change.
Description=Static routes to AWS VPC + tailnet via pro-router (192.168.1.60)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
# `ip route replace` is idempotent — safe to (re)run on every start.
ExecStart=/usr/sbin/ip route replace 10.33.128.0/18 via 192.168.1.60
ExecStart=/usr/sbin/ip route replace 100.64.0.0/10 via 192.168.1.60
# `-` prefix: tolerate the route already being absent on teardown.
ExecStop=-/usr/sbin/ip route del 10.33.128.0/18 via 192.168.1.60
ExecStop=-/usr/sbin/ip route del 100.64.0.0/10 via 192.168.1.60

[Install]
WantedBy=multi-user.target
