# cookbooks/lxc-elasticsearch/files/es-cluster-health.timer
[Unit]
Description=Run es-cluster-health.sh every ~60s

[Timer]
# Type=oneshot self-healing pattern (cf. ~/.claude infrastructure rules):
#   OnBootSec        — first fire after cold boot
#   OnActiveSec      — first fire after `systemctl restart` (boot was long ago)
#   OnUnitInactiveSec — recurring fire, measured from the unit's deactivation
# OnUnitInactiveSec (not OnUnitActiveSec) because the service is Type=oneshot
# with no RemainAfterExit — its active window is ~0, so "after last activation"
# would yield Trigger: n/a.
OnBootSec=60s
OnActiveSec=60s
OnUnitInactiveSec=60s
Unit=es-cluster-health.service

[Install]
WantedBy=timers.target
